# SOURCES/amanda-udp.socket
[Unit]
Description=Amanda Activation Socket

[Socket]
Service=amanda-udp.service
ListenDatagram=10080

[Install]
WantedBy=sockets.target
